# files/common/etc/cron.d/apticron
#
# THIS FILE IS MANAGED BY PUPPET
# /etc/puppet/modules/production/apticron/files/common/etc/cron.d/apticron
#

# apticron
00 7	* * *	root	test -x /usr/sbin/apticron && /usr/sbin/apticron --cron
